# mkosi.extra/etc/systemd/system/getty@tty1.service.d/autologin-eqos.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --skip-login --nonewline --noissue --autologin eqos --noreset --noclear %I $TERM