# mkosi.extra/etc/systemd/system/getty@tty1.service.d/exit-poweroff.conf
[Service]
ExecStopPost=/sbin/poweroff -f